# docker/app/supervisor.conf
[program:ahp-rest]
user=root
directory=/root/ahp_rest/
command=uwsgi --http :8000 --module ahp_rest.wsgi
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
